# docker_conf/nginx_certbot/supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/root/nginx_boot.sh

[program:certbot]
command=/root/certbot_renew.sh
