# manual_conf/web-monitor/web-monitor.timer
[Unit]
Description=web-monitor timer (every 5-10min)

[Timer]
OnCalendar=*:0/10
Persistent=true

[Install]
WantedBy=timers.target

